# deploy/ubuntu/diwan-pen-kiosk.service
[Unit]
Description=Diwan Pen Config Kiosk (User)
After=graphical-session.target network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/diwan-signature
Environment=PYTHONUNBUFFERED=1
Environment=DISPLAY=:0
ExecStart=/opt/diwan-signature/.venv/bin/python -m diwan_signature.desktop.pen_kiosk
Restart=always
RestartSec=2
TimeoutStopSec=20

[Install]
WantedBy=default.target
