# deploy/ubuntu/diwan-signature-flask.service
[Unit]
Description=Diwan Signature Flask Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/diwan-signature
EnvironmentFile=/etc/diwan-signature/diwan-signature.env
ExecStart=/opt/diwan-signature/.venv/bin/python -m diwan_signature serve-flask --host 0.0.0.0 --port 5001
Restart=always
RestartSec=3
KillSignal=SIGINT
TimeoutStopSec=30
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
